# scripts/com.yunoki.message-aggregator-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yunoki.message-aggregator-health</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yunoki/Dropbox/project/message-aggregator/src/message-aggregator/scripts/health-check.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/yunoki/Dropbox/project/message-aggregator/src/message-aggregator</string>

    <key>StandardOutPath</key>
    <string>/Users/yunoki/Dropbox/project/message-aggregator/src/message-aggregator/logs/health-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yunoki/Dropbox/project/message-aggregator/src/message-aggregator/logs/health-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
